# domino/src/handbook-text/tty2.conf
# tty2 - getty
#
# This service maintains a getty on tty1 ssa
# started until it is shut down again.

start on (stopped rc RUNLEVEL=[2375] and started infoserver)
stop on runlevel [!2375]

respawn
exec /bin/login -f handbook TERM=ssa </dev/tty2 >/dev/tty2 2>&1

